# kafka-entity-synthesis/kafka-streamer.service
[Unit]
Description=New Relic Kafka Event Streamer
Documentation=https://github.com/your-org/kafka-entity-synthesis
After=network.target

[Service]
Type=simple
User=newrelic
Group=newrelic
WorkingDirectory=/opt/newrelic/kafka-entity-synthesis
Environment="NODE_ENV=production"
Environment="NODE_OPTIONS=--max-old-space-size=512"

# Main process
ExecStart=/usr/bin/node continuous-exact-format-streamer.js ${CLUSTER_NAME}

# Restart policy
Restart=always
RestartSec=30
StandardOutput=append:/var/log/newrelic/kafka-streamer.log
StandardError=append:/var/log/newrelic/kafka-streamer.error.log

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/newrelic

# Resource limits
LimitNOFILE=65536
MemoryLimit=1G
CPUQuota=50%

[Install]
WantedBy=multi-user.target